# build/systemd/scriptables.service
[Unit]
Description=Scriptables
After=network.target

[Service]
User=www-data
Group=www-data
Restart=always
RestartSec=3
EnvironmentFile=/etc/environment
WorkingDirectory=/var/www/scriptables/
ExecStart=/var/www/scriptables/server
StandardOutput=append:/var/log/scriptables.log
StandardError=append:/var/log/scriptables.error.log

[Install]
WantedBy=multi-user.target